# crontab/cron
*/15 * * * * supervisorctl start gamesrv-backup
5 */2 * * *  supervisorctl start gamesrv-restart
*/5 * * * *  supervisorctl start gamesrv-monitor
